# dot/i3status/checkupdates.service
[Unit]
Description=Check for system updates
Wants=checkupdates.timer

[Service]
Type=oneshot
ExecStart=%h/.config/i3status/checkupdates.sh
